# rootfs/etc/cron.d/owncloud
*/1  *  *  *  * www-data php -f /var/www/owncloud/cron.php
