# rootfs/etc/cron.d/owncloud
*/1  *  *  *  * www-data /usr/bin/cronjob
