# rootfs/etc/cron.d/owncloud
*/1  *  *  *  * www-data /usr/bin/php -f /var/www/owncloud/cron.php
